# usr/lib/systemd/user/jasmine-fixes.service.d/90-jasmine.conf
[Service]
ExecStartPre=setprop vendor.boot.fpsensor $(getprop ro.boot.fpsensor)
ExecStartPre=/bin/sleep 10
ExecStartPre=chmod 777 /sys/class/leds/lcd-backlight/brightness
